# system/systemd/mv-stg-backups/mv-stg-backups.path
[Unit]
Description=Monitor ~/Downloads/librewolf-stg-backups/ for new files and move them/rm the dir
Requires=mv-stg-backups.service

[Path]
PathChanged=/home/austin/Downloads/librewolf-stg-backups/
PathExists=/home/austin/Downloads/librewolf-stg-backups/
Unit=mv-stg-backups.service

[Install]
WantedBy=default.target
