# system/systemd/run-backup/run-backup.service
[Unit]
Description="Service unit (runner) for script that regularly backs up user files"

[Service]
ExecStart=/usr/local/systemd/mv-stg-backups
